# cron/crontab
0 4 * * 1 /etc/cron.d/update.sh
0 0 * * * /etc/cron.d/modif_crontab.sh
@reboot /etc/cron.d/update.sh
